# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/data/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
pidfile=/var/run/supervisord.pid
loglevel=info
user=root

# ── MySQL ──────────────────────────────────────────────────────────────────
[program:mysql]
command=/usr/sbin/mysqld --user=mysql --datadir=/var/lib/mysql --bind-address=0.0.0.0 --character-set-server=utf8mb4 --collation-server=utf8mb4_unicode_ci
directory=/var/lib/mysql
autostart=true
autorestart=true
startretries=3
stdout_logfile=/data/logs/mysql-stdout.log
stderr_logfile=/data/logs/mysql-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
# MySQL 启动较慢，依赖它启动成功后再启动 app
priority=10

# ── Java Spring Boot ───────────────────────────────────────────────────────
[program:app]
command=java -jar /app/app.jar
directory=/app
environment=DB_HOST="localhost",DB_PORT="3306",DB_NAME="tg_downloader",DB_USER="root",DB_PASSWORD="root_password"
autostart=true
autorestart=true
startretries=5
startsecs=15
stdout_logfile=/data/logs/app-stdout.log
stderr_logfile=/data/logs/app-stderr.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
# 等待 MySQL 完全就绪后再启动
priority=20

# ── Nginx ──────────────────────────────────────────────────────────────────
[program:nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf
directory=/usr/share/nginx/html
autostart=true
autorestart=true
startretries=3
stdout_logfile=/data/logs/nginx-stdout.log
stderr_logfile=/data/logs/nginx-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
priority=30
